# config/supervisor/task.ini
[program:pb-client-api]
directory=/go/src/go-ecr-hello
command=/go/bin/hello
stdout_logfile=/tmp/work.api.log
stdout_logfile_maxbytes=50MB
stderr_logfile=/tmp/work.api.err.log
stderr_logfile_maxbytes=50MB
startsecs=5
autorestart=true
startretries=200

[program:pb-client-account]
directory=/go/src/go-ecr-hello
command=/go/bin/user
stdout_logfile=/tmp/work.account.log
stdout_logfile_maxbytes=50MB
stderr_logfile=/tmp/work.account.err.log
stderr_logfile_maxbytes=50MB
startsecs=5
autorestart=true
startretries=200
